# etc/cron.d/bismark-active-tmpfs
* * * * * /tmp/usr/bin/bismark-measure-wrapper
0 */12 * * * /usr/bin/bismark-action scriptupdate
